# supervisor/laravel-websockets.conf
[program:laravel-websockets]
command=php /var/www/html/artisan websockets:serve
directory=/var/www/html
stdout_logfile=/var/www/html/storage/logs/websockets.log
redirect_stderr=true
autostart=true
autorestart=true
user=root